# uwsgi.conf
[program:uwsgi]

command=/home/yevhen/.local/bin/uwsgi students_tracker.wsgi
directory=/home/yevhen/students_tracker/src
autostart=true
autorestart=true
